# patches/sysutils/panicmail/work/.PLIST.flattened
/usr/local/etc/rc.d/panicmail
/usr/local/etc/cperciva-panicmail.pem
/usr/local/share/licenses/panicmail-1.0/catalog.mk
/usr/local/share/licenses/panicmail-1.0/LICENSE
/usr/local/share/licenses/panicmail-1.0/BSD
